# build/etc/supervisor.d/cesi.ini
[program:cesi]
command=/usr/local/cesi/venv/bin/python3 /usr/local/cesi/cesi/run.py --config-file /etc/cesi.conf.toml
#user=root
priority=5
autostart=true
autorestart=true
#startsecs=1
startretries=3
#exitcodes=1,2
stopsignal=QUIT
stopasgroup=false
killasgroup=false
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/cesi.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_events_enabled=false
serverurl=AUTO