# build/etc/supervisor.d/uwsgi.ini
[program:uwsgi]
command=/usr/sbin/uwsgi --ini /etc/cesi/uwsgi.ini
#user=www-data
priority=5
autostart=true
autorestart=true
#startsecs=1
startretries=3
#exitcodes=1,2
stopsignal=QUIT
stopasgroup=false
killasgroup=false
stopwaitsecs=10
redirect_stderr=true
stdout_logfile=/var/log/cesi.log
stdout_logfile_maxbytes=1MB
stdout_logfile_backups=10
stdout_events_enabled=false
serverurl=AUTO
environment=CESI_CONFIG_PATH=/etc/cesi/conf.toml

#process_name=%(program_name)s
#numprocs=1
#directory=/tmp
#umask=022
#stdout_capture_maxbytes=1MB
#stderr_logfile=/a/path
#stderr_logfile_maxbytes=1MB
#stderr_logfile_backups=10
#stderr_capture_maxbytes=1MB
#stderr_events_enabled=false
